# KDE/applications/autostart/ydotoold.desktop
[Desktop Entry]
Exec=ydotoold
Icon=
Name=ydotoold
Path=
Terminal=False
Type=Application
